# rpm/rotate_mongodb.cron
45 * * * * root python /usr/lib/python2.6/site-packages/mongor/tools/rotate.py -s localhost -p 11911 -c filesData -t open_metadata

15 * * * * root python /usr/lib/python2.6/site-packages/mongor/tools/rotate.py -s localhost -p 33933 -c filesData -t sensitive_metadata
